# systemd/mpc-dunst.service
[Unit]
Description = MPC Notification

[Service]
ExecStart = /usr/local/bin/mpc_idle.sh

[Install]
WantedBy = default.target
